# extra/slash/etc/profile.d/libreoffice-still.sh
# to force a certain look'n feel

#export SAL_USE_VCLPLUGIN=gen
#export SAL_USE_VCLPLUGIN=kde4
#export SAL_USE_VCLPLUGIN=gtk
export VCL_ICONS_FOR_DARK_THEME=1
#export SAL_USE_VCLPLUGIN=gtk3
